# setup/apiarium-web.conf
# Ubuntu upstart file at /etc/init/apiarium-web.conf

start on runlevel [2345]
stop on runlevel [016]
 
pre-start script
echo > /var/log/apiarium.log
end script
 
respawn
respawn limit 15 5
 
start on runlevel [2345]
stop on runlevel [06]
 
script
chdir /opt/apiarium
python apiarium web
end script
